# server/game-server.service
[Unit]
Description=Game Server

[Service]
ExecStart=/opt/game-server/start.sh

[Install]
WantedBy=multi-user.target
